# files/router/hotplug.d/iface/30-openvpn
#!/bin/sh
ip route add table openvpn_tun10 default dev tun10
ip route add table openvpn_tun11 default dev tun11
